# sources/meta-ornl/recipes-core/ksz-initscripts/files/init-ksz9897.service
[Unit]
Description=start initscript upon first boot

[Service]
Type=simple
ExecStart=/bin/sh /sbin/init-ksz9897.sh start

[Install]
WantedBy=multi-user.target
